# scripts/com.arroyo.blog-weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Weekly sermon-blog generator. Install:
       cp scripts/com.arroyo.blog-weekly.plist ~/Library/LaunchAgents/
       launchctl load ~/Library/LaunchAgents/com.arroyo.blog-weekly.plist
     Runs Mondays 9:15am. If the Mac is asleep then, launchd runs it on next wake. -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.arroyo.blog-weekly</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/dakotayates/ai-os/ventures/arroyo-church-redesign/scripts/run-blog-weekly.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>1</integer>
    <key>Hour</key><integer>9</integer>
    <key>Minute</key><integer>15</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/dakotayates/Library/Logs/arroyo-blog.launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/dakotayates/Library/Logs/arroyo-blog.launchd.log</string>
</dict>
</plist>
